# 02_atk_mp1/buildroot/buildroot-2020.02.6/package/busybox/inittab
version https://git-lfs.github.com/spec/v1
oid sha256:9bc597c98d8a9e62d08c6357df6bbcf0050d3c898db5c34201f753079e68ba48
size 1241
