# 02_atk_mp1/buildroot/buildroot-2020.02.6/package/sysvinit/inittab
version https://git-lfs.github.com/spec/v1
oid sha256:20adf3ff85c99aab9aa6bd111fdfa9ea959b9c3e07eb23e126506393dfdfe09e
size 1005
