# 02_atk_mp1/buildroot/buildroot-2020.02.6/package/ifupdown-scripts/network/if-pre-up.d/wait_iface
version https://git-lfs.github.com/spec/v1
oid sha256:b70434ea0e10ab723658230c0b52b687e939009886329ee57e4b23008b64c465
size 576
